# new-wallpaper.timer
[Unit]
Description=Periodically change wallpapers.

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
